# src/rootfs/files/configs/etc/systemd/system/check-config-files.timer
[Unit]
Description=Run Check and Copy Config Files Service every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
