# src/rootfs/files/configs/etc/systemd/system/run-state-setup.service
[Unit]
Description=Prepare subdirectories under /run/state for bind mounts
After=run-state.mount
Requires=run-state.mount
Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/run-state-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
