# src/rootfs/files/configs/etc/systemd/system/run-state-setup.service
[Unit]
Description=Prepare /run/state subdirectories for bind mounts
DefaultDependencies=no
After=run-state.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/run-state-setup.sh
RemainAfterExit=yes

[Install]
# Not enabled directly; pulled by mount unit drop-ins
