# src/rootfs/files/configs/etc/systemd/system/sp-svc-workload-service.service
[Unit]
Description=SP Swarm Service - workload-service
After=network-online.target nats-server.service redis-server.service mongod.service
Wants=network-online.target
Requires=nats-server.service redis-server.service mongod.service

[Service]
Type=simple
WorkingDirectory=/usr/local/lib/sp-swarm-services/apps/workload-service
ExecStart=/usr/local/bin/sp-svc-workload-service.sh
Environment=NODE_ENV=production
Restart=always
RestartSec=5
StandardOutput=append:/var/log/sp-svc-workload-service.log
StandardError=append:/var/log/sp-svc-workload-service-err.log

[Install]
WantedBy=multi-user.target
