# src/rootfs/files/configs/etc/systemd/system/sp-tap-network.service
[Unit]
Description=Configure VM network from provider_config spnet settings
DefaultDependencies=no
Before=network-pre.target systemd-networkd.service
Wants=network-pre.target
# Run before networkd reads its config, so the .network file already exists.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/sp-tap-network.sh

[Install]
WantedBy=sysinit.target
