# src/rootfs/files/configs/etc/systemd/system/swarm-cloud-ui.service
[Unit]
Description=Swarm Cloud UI service (static export or Next server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/usr/local/lib/swarm-cloud
ExecStart=/usr/local/bin/swarm-cloud-ui.sh
StandardOutput=append:/var/log/swarm-cloud-ui.log
StandardError=append:/var/log/swarm-cloud-ui-err.log
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
