# src/rootfs/files/configs/pki-service/systemd/pki-authority-sync.service
[Unit]
Description=PKI authority sync client service
After=network-online.target nvidia-persistenced.service pki-vm-mode.service
Wants=network-online.target nvidia-persistenced.service pki-vm-mode.service
Requires=pki-vm-mode.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=infinity
Environment=NODE_ENV=production
ExecCondition=/bin/bash -c '[[ "$(head -n1 /etc/swarm/swarm-vm-mode 2>/dev/null | tr -d "[:space:]")" == "normal" ]]'
ExecStartPre=/usr/bin/python3 /usr/local/bin/pki-authority/pki_configure_helper.py \
    configure-sync-client \
    --config /sp/swarm/config.yaml \
    --output /etc/swarm/sync-client-config.yaml
ExecStart=/bin/bash -c 'until /usr/bin/node /usr/bin/pki-sync-client sync \
    --config /etc/super/pki-authority-sync/secrets-config.yaml \
    --swarm-env /etc/swarm/sync-client-config.yaml \
    --save-certs-to-dir /etc/super/certs/vm \
    --save-certs-prefix vm \
    --disable-server-identity-check \
    --log-level debug; do \
    echo "[pki-authority-sync] remote not ready, retry in 60s"; \
    sleep 60; \
done'

[Install]
WantedBy=multi-user.target
