# src/rootfs/files/configs/pki-service/systemd/pki-cert-init.service
[Unit]
Description=PKI certificate initialization service
After=nvidia-persistenced.service pccs.service pki-vm-mode.service
Wants=nvidia-persistenced.service pki-vm-mode.service
Requires=pccs.service pki-vm-mode.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=infinity
ExecCondition=/bin/bash -c '[[ "$(head -n1 /etc/swarm/swarm-vm-mode 2>/dev/null | tr -d "[:space:]")" == "init" ]]'
ExecStartPre=/usr/bin/python3 /usr/local/bin/pki-authority/pki_cert_gen_configure.py \
    configure \
    --template /etc/super/pki-cert-generator/cert-gen-config-template.yaml \
    --output /etc/super/pki-cert-generator/cert-gen-config.yaml
ExecStart=/bin/bash -c 'until /usr/bin/pki-cert-generator generate-chain \
    -c /etc/super/pki-cert-generator/cert-gen-config.yaml \
    --log-level debug; do \
    echo "[pki-cert-init] generation failed, retry in 60s"; \
    sleep 60; \
done'

[Install]
WantedBy=multi-user.target