# src/rootfs/files/configs/pki-service/systemd/pki-authority-watchdog.timer
[Unit]
Description=PKI Authority Watchdog Timer


[Timer]
OnUnitActiveSec=5min
Unit=pki-authority-watchdog.service

[Install]
WantedBy=timers.target
