# src/rootfs/files/configs/pki-service/systemd/pki-authority-watchdog.service
[Unit]
Description=PKI Authority Watchdog

[Service]
Type=oneshot
TimeoutStartSec=1min
ExecStart=/usr/local/bin/pki-authority/pki_watchdog.sh
