# src/rootfs/files/configs/state_disk_mount/state_disk_mount.service
[Unit]
Description=Create LUKS partition with session key
DefaultDependencies=no
Before=local-fs.target cryptsetup.target
After=systemd-udevd.service
Conflicts=shutdown.target
## Do not pin to a specific device name; disks may enumerate differently on some hosts.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/state_disk_mount.sh
TimeoutStartSec=300

[Install]
WantedBy=local-fs.target
