# src/rootfs/files/configs/pki-service/systemd/pki-authority.service.d/init.conf
[Unit]
# Auto-start only in swarm-init mode
ConditionKernelCommandLine=vm_mode=swarm-init

[Install]
WantedBy=multi-user.target
